# 30-sentinel.rules
ACTION=="add",DEVPATH=="/module/wiegand_gpio",SUBSYSTEM=="module",RUN+="/bin/chgrp -R sentinel /sys/kernel/wiegand"

